# forgejo/forgejo-dump.service
[Unit]
Description=Forgejo nightly dump + offsite copy to aura-farming
After=forgejo.service
Wants=forgejo.service

[Service]
Type=oneshot
ExecStart=%h/repos/devcontainer-base/forgejo/forgejo-backup.sh
