# forgejo/forgejo-dump.timer
[Unit]
Description=Nightly Forgejo dump

[Timer]
# 03:30 local time daily; Persistent catches up a missed run after the laptop was asleep.
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
